# supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=/usr/sbin/nginx

[program:uwsgi]
user=quilt
command=uwsgi --ini /etc/uwsgi.ini
